# deploy/backend.service.example
[Unit]
Description=Management Assistant Backend Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/management-system/apps/backend
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=4000
Environment=TRUST_PROXY=loopback
Environment=ACCESS_TOKEN_SECRET=change-me-access
Environment=REFRESH_TOKEN_SECRET=change-me-refresh
Environment=ACCOUNT_ARCHIVE_SECRET=change-me-archive
Environment=COOKIE_SECURE=true
Environment=ALLOWED_ORIGINS=https://your-domain.example
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
